# brawler.service
[Unit]
Description=Brawler Algo Trading Strategy (Headless)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/slipstream
EnvironmentFile=/home/ubuntu/slipstream/.env.brawler
ExecStart=/home/ubuntu/slipstream/.venv/bin/python -m slipstream.strategies.brawler.cli --config /home/ubuntu/slipstream/config/brawler_eth.yaml
Restart=always
RestartSec=5
StandardOutput=append:/home/ubuntu/slipstream/logs/brawler.log
StandardError=append:/home/ubuntu/slipstream/logs/brawler.err.log

[Install]
WantedBy=multi-user.target
