# brawler.service
[Unit]
Description=Brawler Market Making Strategy
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/slipstream
ExecStart=/home/ubuntu/slipstream/.venv/bin/python -m slipstream.strategies.brawler.cli --config prod_brawler.yaml
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
