# brawler.service
[Unit]
Description=Brawler Market Maker Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/slipstream
EnvironmentFile=/home/ubuntu/slipstream/.env.brawler
Environment=PYTHONUNBUFFERED=1
# Run python directly to ensure singular process management
ExecStart=/home/ubuntu/.local/bin/uv run python -m slipstream.strategies.brawler.cli --config config/brawler_prod.yaml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
